# dot_bash/crontab/crontab.editor
SHELL=/bin/bash
MAILTO=""
0 8 * * * eval 'cleanup(){ TRASH_FILES=( '\''.DS_Store'\'' '\''._.DS_Store'\''); find /Users/ycorrales -type f \( -name ${TRASH_FILES[0]} -o -name ${TRASH_FILES[1]} \) -ls -delete 2> /dev/null; return 0; }; cleanup > /dev/null'
